# systemd/intcheck.timer
[Timer]
OnCalendar=*:*:00/10
Unit=internetcheck.service

