# systemd/intcheck.service
[Service]
ExecStart=/root/system-fd/ping.sh

